# Notification Auto/systemd/glpi-cron.service
[Unit]
Description=GLPI actions automatique

[Service]
User=www-data
Group=www-data
ExecStart=/usr/bin/php /var/www/html/front/cron.php
