# Notification Auto/systemd/glpi-cron.timer
[Unit]
Description=Run GLPI cron every minute

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true

[Install]
WantedBy=timers.target
